# queue/schedules/nbiai-queue-processor.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Windows Task Scheduler XML config for NBIAI Queue Processor.

  Import with:
    schtasks /create /xml "D:\OneDrive\Claude_code\NBIAI_TEAM\queue\schedules\nbiai-queue-processor.xml" /tn "NBIAI\Queue Processor"

  Runs every 15 minutes during working hours (07:00-22:00) on weekdays.
  Only runs when the user is logged in (needs clipboard access for prompt copy).
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Process NBIAI task queue inbox - picks oldest task, copies prompt to clipboard, opens Claude Desktop</Description>
    <Author>NBI Analytics</Author>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <Repetition>
        <Interval>PT15M</Interval>
        <Duration>PT15H</Duration>
        <StopAtDurationEnd>false</StopAtDurationEnd>
      </Repetition>
      <StartBoundary>2026-03-28T07:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByWeek>
        <DaysOfWeek>
          <Monday />
          <Tuesday />
          <Wednesday />
          <Thursday />
          <Friday />
        </DaysOfWeek>
        <WeeksInterval>1</WeeksInterval>
      </ScheduleByWeek>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <ExecutionTimeLimit>PT5M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-ExecutionPolicy Bypass -File "D:\OneDrive\Claude_code\NBIAI_TEAM\queue\scripts\process-inbox.ps1"</Arguments>
      <WorkingDirectory>D:\OneDrive\Claude_code\NBIAI_TEAM</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
